# systemd/lights-out-daemon.service
[Unit]
Description=MSI CORELIQUID temperature monitor (sends CPU temp to cooler)

[Service]
Type=simple
ExecStart=/home/ben/src/lights-out/lights-out/target/release/lights-out daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target
